# DeskThingServer/resources/scripts/restart_script/rootfs/etc/udev/rules.d/50-usb.rules
ACTION=="change", SUBSYSTEM=="android_usb", KERNEL=="android0", RUN+="/sbin/restart_usb"
